# blinkstick-net.service
[Unit]
Description= Blinkstick Network Visualizer

[Install]
WantedBy=default.target

[Service]
ExecStart=/usr/local/bin/blinkstick-net.sh

