# rgb-sound.service
[Unit]
Description= RGB Sound Visualizer

[Install]
WantedBy=default.target

[Service]
ExecStart=/usr/local/bin/rgb-sound.sh

