# rgb-sound-visualizer.service
[Unit]
Description= RGB Sound Visualizer

[Install]
WantedBy=basic.target

[Service]
ExecStart=/usr/local/bin/rgb-sound-visualizer.sh

